# battery-notify.timer
[Unit]
Description=Timer to run battery-notify

[Timer]
OnActiveSec=1min
OnUnitActiveSec=5s
AccuracySec=1us

[Install]
WantedBy=timers.target
