# battery-notify.service
[Unit]
Description=battery-notify

[Service]
ExecStart=%h/bin/check-batt

[Install]
WantedBy=default.target
